# podman-quadlet/pushover-webhook.container
[Unit]
Description=Jellyfin Pushover Webhook Quadlet

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target

[Container]
Image=localhost/pushover-webhook:latest
ContainerName=pushover-webhook
HostName=pushover-webhook

PublishPort=8484:8484

Environment=AUTH_TOKEN=
Environment=PUSHOVER_API_TOKEN=
Environment=PUSHOVER_USER_KEY=
